# projectroot/usr/lib/systemd/system/rauc-hawkbit-updater.service.d/20-install-initial-config.conf
[Service]
ExecStartPre=+/usr/bin/mkdir -p /config/rauc-hawkbit-updater
ExecStartPre=+/usr/bin/cp -n /etc/rauc-hawkbit-updater/config.conf.initial /config/rauc-hawkbit-updater/config.conf
