# local_src/juconn-20220315/etc/NetworkManager/dispatcher.d/10-restart-timesyncd
#!/bin/sh

IF=$1
ACTION=$2

if [ "$ACTION" = "up" ]; then
        systemctl restart systemd-timesyncd
fi

